# docker/files/chatmail-certmon.timer
[Unit]
Description=Periodically check TLS certificate changes

[Timer]
OnBootSec=120
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
