# docker/files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; acmetool-redirector.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:acmetool-redirector]
command=/usr/bin/acmetool redirector --service.uid=daemon
autostart=false
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/acmetool-redirector.log
stderr_logfile=/var/log/supervisor/acmetool-redirector.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; chatmail-metadata.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:chatmail-metadata]
command=/usr/local/lib/chatmaild/venv/bin/chatmail-metadata /run/chatmail-metadata/metadata.socket /usr/local/lib/chatmaild/chatmail.ini
user=vmail
umask=0077
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/chatmail-metadata.log
stderr_logfile=/var/log/supervisor/chatmail-metadata.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; cron.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:cron]
command=/usr/sbin/cron -f $EXTRA_OPTS
environment=EXTRA_OPTS=""
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cron.log
stderr_logfile=/var/log/supervisor/cron.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; dbus.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork --nopidfile
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; doveauth.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:doveauth]
command=/usr/local/lib/chatmaild/venv/bin/doveauth /run/doveauth/doveauth.socket /usr/local/lib/chatmaild/chatmail.ini
user=vmail
umask=0077
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/doveauth.log
stderr_logfile=/var/log/supervisor/doveauth.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; dovecot.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:dovecot]
command=/usr/sbin/dovecot -F
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/dovecot.log
stderr_logfile=/var/log/supervisor/dovecot.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; echobot.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:echobot]
command=/usr/local/lib/chatmaild/venv/bin/echobot /usr/local/lib/chatmaild/chatmail.ini
user=echobot
directory=/var/lib/echobot
umask=0077
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/echobot.log
stderr_logfile=/var/log/supervisor/echobot.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; fcgiwrap.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:fcgiwrap]
command=/usr/sbin/fcgiwrap -f
user=www-data
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/fcgiwrap.log
stderr_logfile=/var/log/supervisor/fcgiwrap.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; filtermail-incoming.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:filtermail-incoming]
command=/usr/local/lib/chatmaild/venv/bin/filtermail /usr/local/lib/chatmaild/chatmail.ini incoming
user=vmail
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/filtermail-incoming.log
stderr_logfile=/var/log/supervisor/filtermail-incoming.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; filtermail.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:filtermail]
command=/usr/local/lib/chatmaild/venv/bin/filtermail /usr/local/lib/chatmaild/chatmail.ini outgoing
user=vmail
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/filtermail.log
stderr_logfile=/var/log/supervisor/filtermail.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; iroh-relay.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:iroh-relay]
command=/usr/local/bin/iroh-relay --config-path /etc/iroh-relay.toml
user=iroh
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/iroh-relay.log
stderr_logfile=/var/log/supervisor/iroh-relay.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; lastlogin.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:lastlogin]
command=/usr/local/lib/chatmaild/venv/bin/lastlogin /run/chatmail-lastlogin/lastlogin.socket /usr/local/lib/chatmaild/chatmail.ini
user=vmail
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/lastlogin.log
stderr_logfile=/var/log/supervisor/lastlogin.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; nginx.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:nginx]
command=/usr/sbin/nginx -g "daemon off; master_process on;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; opendkim.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:opendkim]
command=/usr/sbin/opendkim -f
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/opendkim.log
stderr_logfile=/var/log/supervisor/opendkim.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; postfix@-.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:postfix]
command=/usr/sbin/postfix start-fg
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postfix.log
stderr_logfile=/var/log/supervisor/postfix.err

;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
; unbound.service
;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;;
[program:unbound]
command=/usr/sbin/unbound -d
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/unbound.log
stderr_logfile=/var/log/supervisor/unbound.err
